# sbx-h6-rgb-git/99-sbx-h6-systemd.rules
ACTION=="add", ATTR{idVendor}=="041e", ATTR{idProduct}=="3255", RUN+="/usr/bin/systemctl --no-block restart sbx-h6.service"
